# runWaffleServer
#!/bin/sh /etc/rc.common

START=80
STOP=15

start() {        
        echo start
        python /root/wafflecarServer.py
}                 

stop() {          
        echo stop
}
